# cisa-alerts/aa20-336a.json
{
  "markdown": "\n\nSummary\n\nThis Advisory uses the MITRE Adversarial Tactics, Techniques, and Common Knowledge (ATT&CK\u00ae) framework. See the ATT&CK for Enterprise for all referenced threat actor tactics and techniques.\nThe Cybersecurity and Infrastructure Security Agency (CISA) and the Federal Bureau of Investigation (FBI) have observed persistent continued cyber intrusions by advanced persistent threat (APT) actors targeting U.S. think tanks. This malicious activity is often, but not exclusively, directed at individuals and organizations that focus on international affairs or national security policy.[1] The following guidance may assist U.S. think tanks in developing network defense procedures to prevent or rapidly detect these attacks.\nAPT actors have relied on multiple avenues for initial access. These have included low-effort capabilities such as spearphishing emails and third-party message services directed at both corporate and personal accounts, as well as exploiting vulnerable web-facing devices and remote connection capabilities. Increased telework during the COVID-19 pandemic has expanded workforce reliance on remote connectivity, affording malicious actors more opportunities to exploit those connections and to blend in with increased traffic. Attackers may leverage virtual private networks (VPNs) and other remote work tools to gain initial access or persistence on a victim\u2019s network. When successful, these low-effort, high-reward approaches allow threat actors to steal sensitive information, acquire user credentials, and gain persistent access to victim networks.\nGiven the importance that think tanks can have in shaping U.S. policy, CISA and FBI urge individuals and organizations in the international affairs and national security sectors to immediately adopt a heightened state of awareness and implement the critical steps listed in the Mitigations section of this Advisory.\nClick here for a PDF version of this report.\n\nTechnical Details\n\nATT&CK Profile\nCISA created the following MITRE ATT&CK profile to provide a non-exhaustive list of tactics, techniques, and procedures (TTPs) employed by APT actors to break through think tanks\u2019 defenses, conduct reconnaissance in their environments, exfiltrate proprietary or confidential information, and execute effects on targets. These TTPs were included based upon closed reporting on APT actors that are known to target think tanks or based upon CISA incident response data.\n\nInitial Access [TA0001]\n\nValid Accounts [T1078]\nValid Accounts: Cloud Accounts [T1078.004]\nExternal Remote Services [T1133]\nDrive-by Compromise [T1189]\nExploit Public-Facing Application [T1190]\n\nSupply Chain Compromise: Compromise Software Supply Chain [T1195.002]\nTrusted Relationship [T1199]\nPhishing: Spearphishing Attachment [T1566.001]\nPhishing: Spearphishing Link [T1566.002]\nPhishing: Spearphishing via Service [T1566.003]\n\n\n\n\nExecution [TA0002]\n\nWindows Management Instrumentation [T1047]\nScheduled Task/Job: Scheduled Task [T1053.005]\nCommand and Scripting Interpreter: PowerShell [T1059.001]\nCommand and Scripting Interpreter: Windows Command Shell [T1059.003]\nCommand and Scripting Interpreter: Unix Shell [T1059.004]\nCommand and Scripting Interpreter: Visual Basic [T1059.005]\nCommand and Scripting Interpreter: Python [T1059.006]\nNative API [T1106]\nExploitation for Client Execution [T1203]\nUser Execution: Malicious Link [T1204.001]\nUser Execution: Malicious File [T1204.002]\nInter-Process Communication: Dynamic Data Exchange [T1559.002]\nSystem Services: Service Execution [T1569.002]\n\n\nPersistence [TA0003]\n\nBoot or Logon Initialization Scripts: Logon Script (Windows) [T1037.001]\nScheduled Task/Job: Scheduled Task [T1053.005]\nAccount Manipulation: Exchange Email Delegate Permissions [T1098.002]\nCreate Account: Local Account [T1136.001]\nOffice Application Startup: Office Test [T1137.002]\nOffice Application Startup: Outlook Home Page [T1137.004]\nBrowser Extensions [T1176]\nBITS Jobs [T1197]\nServer Software Component: Web Shell [T1505.003]\nPre-OS Boot: Bootkit [T1542.003]\nCreate or Modify System Process: Windows Service [T1543.003]\nEvent Triggered Execution: Change Default File Association [T1546.001]\nEvent Triggered Execution: Windows Management Instrumentation Event Subscription [T1546.003]\nEvent Triggered Execution: Accessibility Features [T1546.008]\nEvent Triggered Execution: Component Object Model Hijacking [T1546.015]\nBoot or Logon Autostart Execution: Registry Run Keys / Startup Folder [T1547.001]\nBoot or Logon Autostart Execution: Shortcut Modification [T1547.009]\n\n\nPrivilege Escalation [TA0004]\n\nProcess Injection [T1055]\nProcess Injection: Process Hollowing [T1055.012]\nExploitation for Privilege Escalation [T1068]\nAccess Token Manipulation: Token Impersonation/Theft [T1134.001]\nEvent Triggered Execution: Accessibility Features [T1546.008]\nBoot or Logon Autostart Execution: Shortcut Modification [T1547.009]\nAbuse Elevation Control Mechanism: Bypass User Access Control [T1548.002]\nHijack Execution Flow: DLL Side-Loading [T1574.002]\n\n\nDefense Evasion [TA0005]\n\nRootkit [T1014]\nObfuscated Files or Information: Binary Padding [T1027.001]\nObfuscated Files or Information: Software Packing [T1027.002]\nObfuscated Files or Information: Steganography [T1027.003]\nObfuscated Files or Information: Indicator Removal from Tools [T1027.005]\nMasquerading: Match Legitimate Name or Location [T1036.005]\nIndicator Removal on Host: Clear Windows Event Logs [T1070.001]\nIndicator Removal on Host: Clear Command History [1070.003]\nIndicator Removal on Host: File Deletion [T1070.004]\nIndicator Removal on Host: Timestomp [T1070.006]\nModify Registry [T1112]\nDeobfuscate/Decode Files or Information [T1140]\nExploitation for Defense Evasion [T1211]\nSigned Binary Proxy Execution: Compiled HTML File [T1218.001]\nSigned Binary Proxy Execution: Mshta [T1218.005]\nSigned Binary Proxy Execution: Rundll32 [T1218.011]\nTemplate Injection [T1221]\nExecution Guardrails: Environmental Keying [T1480.001]\nAbuse Elevation Control Mechanism: Bypass User Access Control [T1548.002]\nUse Alternate Authentication Material: Application Access Token [T1550.001]\nSubvert Trust Controls: Code Signing [T1553.002]\nImpair Defenses: Disable or Modify Tools [T1562.001]\nImpair Defenses: Disable or Modify System Firewall [T1562.004]\nHide Artifacts: Hidden Files and Directories [T1564.001]\nHide Artifacts: Hidden Window [T1564.003]\n\n\nCredential Access [TA0006]\n\nOS Credential Dumping: LSASS Memory [T1003.001]\nOS Credential Dumping: Security Account Manager [T1003.002]\nOS Credential Dumping: NTDS [T1003.003]\nOS Credential Dumping: LSA Secrets [T1003.004]\nOS Credential Dumping: Cached Domain Credentials [T1003.005]\nNetwork Sniffing [T1040]\nInput Capture: Keylogging [T1056.001]\nBrute Force: Password Cracking [T1110.002]Brute Force: Password Spraying [T1110.003]\nForced Authentication [T1187]\nSteal Application Access Token [T1528]\nUnsecured Credentials: Credentials in Files [T1552.001]\nUnsecured Credentials: Group Policy Preferences [T1552.006]\nCredentials from Password Stores: Credentials from Web Browsers [T1555.003]\n\n\nDiscovery [TA0007]\n\nSystem Service Discovery [T1007]\nQuery Registry [T1012]\nSystem Network Configuration Discovery [T1016]\nRemote System Discovery [T1018]\nSystem Owner/User Discovery [T1033]\nNetwork Sniffing [T1040]\nNetwork Service Scanning [T1046]\nSystem Network Connections Discovery [T1049]\nProcess Discovery [T1057]\nPermission Groups Discovery: Local Groups [T1069.001]\nPermission Groups Discovery: Domain Groups [T1069.002]\nSystem Information Discovery [T1082]\nFile and Directory Discovery [T1083]\nAccount Discovery: Local Account [T1087.001]\nAccount Discovery: Domain Account [T1087.002]\nPeripheral Device Discovery [T1120]\nNetwork Share Discovery [T1135]\nPassword Policy Discovery [T1201]\nSoftware Discovery: Security Software Discovery [T1518.001]\n\n\nLateral Movement [TA0008]\n\nRemote Services: Remote Desktop Protocol [T1021.001]\nRemote Services: SSH [T1021.004]\nTaint Shared Content [T1080]\nReplication Through Removable Media [T1091]\nExploitation of Remote Services [T1210]\nUse Alternate Authentication Material: Pass the Hash [T1550.002]\nUse Alternate Authentication Material: Pass the Ticket [T1550.003]\n\n\nCollection [TA0009]\n\nData from Local System [T1005]\nData from Removable Media [T1025]\nData Staged: Local Data Staging [T1074.001]\nScreen Capture [T1113]\nEmail Collection: Local Email Collection [T1114.001]\nEmail Collection: Remote Email Collection [T1114.002]\nAutomated Collection [T1119]\nAudio Capture [T1123]\nData from Information Repositories: SharePoint [T1213.002]\nArchive Collected Data: Archive via Utility [T1560.001]\nArchive Collected Data: Archive via Custom Method [T1560.003]\n\n\nCommand and Control [TA0011]\n\nData Obfuscation: Junk Data [T1001.001]\nFallback Channels [T1008]\nApplication Layer Protocol: Web Protocols [T1071.001]\nApplication Layer Protocol: File Transfer Protocols [T1071.002]\nApplication Layer Protocol: Mail Protocols [T1071.003]\nApplication Layer Protocol: DNS [T1071.004]\nProxy: External Proxy [T1090.002]\nProxy: Multi-hop Proxy [T1090.003]\nProxy: Domain Fronting [T1090.004]\nCommunication Through Removable Media [T1092]\nNon-Application Layer Protocol [T1095]\nWeb Service: Dead Drop Resolver [T1102.001]\nWeb Service: Bidirectional Communication [T1102.002]\nMulti-Stage Channels [T1104]\nIngress Tool Transfer [T1105]\nData Encoding: Standard Encoding [T1132.001]\nRemote Access Software [T1219]\nDynamic Resolution: Domain Generation Algorithms [T1568.002]\nNon-Standard Port [T1571]\nProtocol Tunneling [T1572]\nEncrypted Channel: Symmetric Cryptography [T1573.001]\nEncrypted Channel: Asymmetric Cryptography [T1573.002]\n\n\n\u00a0Exfiltration [TA0010]\n\nExfiltration Over C2 Channel [T1041]\nExfiltration Over Alternative Protocol: Exfiltration Over Unencrypted/Obfuscated Non-C2 Protocol [T1048.003]\n\n\nImpact [TA0040]\n\nData Encrypted for Impact [T1486]\nResource Hijacking [T1496]\nSystem Shutdown/Reboot [T1529]\nDisk Wipe: Disk Structure Wipe [T1561.002]\n\n\n\n\nMitigations\n\nCISA and FBI recommend think tank organizations apply the following critical practices to strengthen their security posture.\nLeaders\n\nImplement a training program to familiarize users with identifying social engineering techniques and phishing emails.\n\nUsers/Staff\n\nLog off remote connections when not in use.\nBe vigilant against tailored spearphishing attacks targeting corporate and personal accounts (including both email and social media accounts).\nUse different passwords for corporate and personal accounts.\nInstall antivirus software on personal devices to automatically scan and quarantine suspicious files.\nEmploy strong multi-factor authentication for personal accounts, if available.\nExercise caution when:\n\nOpening email attachments, even if the attachment is expected and the sender appears to be known. See Using Caution with Email Attachments.\nUsing removable media (e.g., USB thumb drives, external drives, CDs).\n\n\n\nIT Staff/Cybersecurity Personnel\n\nSegment and segregate networks and functions.\nChange the default username and password of applications and appliances.\nEmploy strong multi-factor authentication for corporate accounts.\nDeploy antivirus software on organizational devices to automatically scan and quarantine suspicious files.\nApply encryption to data at rest and data in transit.\nUse email security appliances to scan and remove malicious email attachments or links.\nMonitor key internal security tools and identify anomalous behavior. Flag any known indicators of compromise or threat actor behaviors for immediate response.\nOrganizations can implement mitigations of varying complexity and restrictiveness to reduce the risk posed by threat actors who use Tor (The Onion Router) to carry out malicious activities. See the CISA-FBI Joint Cybersecurity Advisory on Defending Against Malicious Cyber Activity Originating from Tor for mitigation options and additional information.\nPrevent exploitation of known software vulnerabilities by routinely applying software patches and upgrades. Foreign cyber threat actors continue to exploit publicly known\u2014and often dated\u2014software vulnerabilities against broad target sets, including public and private sector organizations. If these vulnerabilities are left unpatched, exploitation often requires few resources and provides threat actors with easy access to victim networks. Review CISA and FBI\u2019s Top 10 Routinely Exploited Vulnerabilities and other CISA alerts that identify vulnerabilities exploited by foreign attackers.\nImplement an antivirus program and a formalized patch management process.\nBlock certain websites and email attachments commonly associated with malware (e.g., .scr, .pif, .cpl, .dll, .exe).\nBlock email attachments that cannot be scanned by antivirus software (e.g., .zip files).\nImplement Group Policy Object and firewall rules.\nImplement filters at the email gateway and block suspicious IP addresses at the firewall.\nRoutinely audit domain and local accounts as well as their permission levels to look for situations that could allow an adversary to gain wide access by obtaining credentials of a privileged account.\nFollow best practices for design and administration of the network to limit privileged account use across administrative tiers.\nImplement a Domain-Based Message Authentication, Reporting & Conformance (DMARC) validation system.\nDisable or block unnecessary remote services.\nLimit access to remote services through centrally managed concentrators.\nDeny direct remote access to internal systems or resources by using network proxies, gateways, and firewalls.\nLimit unnecessary lateral communications.\nDisable file and printer sharing services. If these services are required, use strong passwords or Active Directory authentication.\nEnsure applications do not store sensitive data or credentials insecurely.\nEnable a firewall on agency workstations, configured to deny unsolicited connection requests.\nDisable unnecessary services on agency workstations and servers.\nScan for and remove suspicious email attachments; ensure any scanned attachment is its \"true file type\" (i.e., the extension matches the file header).\nMonitor users' web browsing habits; restrict access to suspicious or risky sites. Contact law enforcement or CISA immediately regarding any unauthorized network access identified.\nVisit the MITRE ATT&CK techniques and tactics pages linked in the ATT&CK Profile section above for additional mitigation and detection strategies for this malicious activity targeting think tanks.\n\n\nContact Information\n\nRecipients of this report are encouraged to contribute any additional information that they may have related to this threat. To report suspicious or criminal activity related to information found in this Joint Cybersecurity Advisory, contact your local FBI field office at www.fbi.gov/contact-us/field, or the FBI\u2019s 24/7 Cyber Watch (CyWatch) at (855) 292-3937 or by email at CyWatch@fbi.gov. When available, please include the following information regarding the incident: date, time, and location of the incident; type of activity; number of people affected; type of equipment used for the activity; the name of the submitting company or organization; and a designated point of contact. To request incident response resources or technical assistance related to these threats, contact CISA at Central@cisa.gov.\nReferences\n\nCISA Alert: Microsoft Office 365 Security Recommendations\nCISA Alert: Technical Approaches to Uncovering and Remediating Malicious Activity\nCISA Webpage: Telework Guidance\nCISA Webpage: VPN-Related Guidance\nFBI Private Industry Notification: PIN 20200409-001\n\n\nReferences\n\n[1] CyberScoop: As Europe prepares to vote, Microsoft warns of Fancy Bear attacks on democratic think tanks\n\nRevisions\n\nInitial Version: December 1, 2020\n",
  "cves": [],
  "techniques": [
    "T1505.003",
    "T1134.001",
    "T1070.004",
    "T1560.001",
    "T1021.004",
    "T1136.001",
    "T1053.005",
    "T1059.004",
    "T1574.002",
    "T1566.001",
    "T1568.002",
    "T1561.002",
    "T1003.003",
    "T1546.003",
    "T1059.001",
    "T1550.002",
    "T1137.002",
    "T1071.001",
    "T1027.003",
    "T1090.003",
    "T1132.001",
    "T1003.005",
    "T1204.001",
    "T1059.003",
    "T1562.004",
    "T1071.002",
    "T1204.002",
    "T1059.005",
    "T1110.002",
    "T1552.001",
    "T1480.001",
    "T1110.003",
    "T1071.003",
    "T1548.002",
    "T1070.006",
    "T1573.002",
    "T1098.002",
    "T1027.002",
    "T1069.001",
    "T1569.002",
    "T1518.001",
    "T1546.015",
    "T1001.001",
    "T1087.002",
    "T1559.002",
    "T1555.003",
    "T1003.004",
    "T1546.008",
    "T1027.005",
    "T1218.001",
    "T1102.001",
    "T1114.002",
    "T1137.004",
    "T1218.005",
    "T1560.003",
    "T1048.003",
    "T1547.001",
    "T1550.001",
    "T1573.001",
    "T1564.003",
    "T1027.001",
    "T1021.001",
    "T1547.009",
    "T1546.001",
    "T1074.001",
    "T1078.004",
    "T1087.001",
    "T1003.001",
    "T1566.003",
    "T1552.006",
    "T1036.005",
    "T1543.003",
    "T1553.002",
    "T1550.003",
    "T1542.003",
    "T1003.002",
    "T1102.002",
    "T1056.001",
    "T1037.001",
    "T1070.001",
    "T1090.002",
    "T1071.004",
    "T1566.002",
    "T1055.012",
    "T1218.011",
    "T1069.002",
    "T1090.004",
    "T1195.002",
    "T1059.006",
    "T1114.001",
    "T1562.001",
    "T1213.002",
    "T1564.001"
  ],
  "advisory": "aa20-336a",
  "title": "Advanced Persistent Threat Actors Targeting U.S. Think Tanks | CISA",
  "source": "cybersecurity-advisories",
  "id": "3807430c32364f3b5d8900fb631d86032b4257916dc14d694748e3dde0185af2"
}